# frontend/launchpad-frontend.service
[Unit]
Description=Launchpad Frontend (Vite preview)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/launchpad/frontend
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run serve -- --host 0.0.0.0 --port 5173
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
